# Backend/WorkerA/scripts/mg-worker-a.service
[Unit]
Description=MG Service (A)
After=network.target

[Service]
User=ubuntu
ExecStart=/home/ubuntu/app/scripts/start.sh

[Install]
WantedBy=multi-user.target